# files/every-offlineimap.service
[Unit]
Description=Start checking all IMAP mailboxes

[Service]
Type=oneshot
ExecStart=/bin/true
